# file transfer from Pi/rc.local
#!/bin/sh -e

sudo /home/pi/on_start.sh & > /home/pi/log.txt
exit 0